# docker/supervisor/conf.d/rabbitmq-worker.conf
[program:rabbitmq-worker]
command=drush rabbitmq:worker audit_log --rabbitmq_timeout=3600
autorestart=true
user=wodby
