# docker/supervisor/conf.d/rabbitmq-worker-queue_mail.conf
[program:rabbitmq-worker-queue_mail]
command=drush rabbitmq:worker queue_mail --rabbitmq_timeout=3600
autorestart=true
user=wodby
